# .docker/oracle.ini
[program:oracle]
command=node src/app3.js
directory=/var/www/oracle-lite
user=root
autostart=true
autorestart=true

[program:haproxy]
command=haproxy -f /etc/haproxy/haproxy.cfg
user=root
autostart=true
autorestart=true